# build/openvoxview.service
[Unit]
]Description=PuppetDB Dashboard
After=network.target
Documentation=https://github.com/voxpupuli/openvoxview

[Service]
Environment=GIN_MODE=release
ExecStart=/usr/bin/openvoxview --config /etc/voxpupuli/openvoxview.yml
ExecReload=/bin/kill -HUP $PID
KillMode=process

[Install]
WantedBy=multi-user.target